# scripts/pve-host-config-backup-to-pbs.timer
[Unit]
Description=Daily PVE host config backup to PBS

[Timer]
OnCalendar=*-*-* 03:45:00
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
